# supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:backend]
command=node server.js
directory=/app/backend
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
environment=MONGO_URI="%(ENV_MONGO_URI)s",JWT_SECRET="%(ENV_JWT_SECRET)s",INTERNAL_API_KEY="%(ENV_INTERNAL_API_KEY)s"

[program:trigger-engine]
command=node index.js
directory=/app/trigger-engine
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
environment=MONGO_URI="%(ENV_MONGO_URI)s",INTERNAL_API_KEY="%(ENV_INTERNAL_API_KEY)s"

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
